# cmdb/deploy/cmdb.ini
[program:cmdb-web]
command=/root/.pyenv/shims/uwsgi --ini cmdb_uwsgi.ini
directory=/data/sinochem/cmdb/deploy/
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cmdb/cloud-uwsgi.log
loglevel=info


[program:cmdb-scheduler-celery]
environment=PYTHONOPTIMIZE=1
command=/root/.pyenv/shims/celery -A cmdb worker --loglevel=info --concurrency=4 -Q cmdb --workdir=/data/sinochem/cmdb
directory=/data/sinochem/cmdb/
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cmdb/cmdb-scheduler-celery.log
loglevel=info